# postgres/backups/etc/periodic/hourly/backup
#!/bin/sh

# Important! Script have to be without extension!

/usr/bin/coverme backup -c /var/backups/config.yml
